# appendix-b-docker/files/etc/s6-overlay/s6-rc.d/cron-log-prepare/up
#!/bin/sh
if { mkdir -p /var/log/cron }
if { chown nobody:nogroup /var/log/cron }
chmod 02755 /var/log/cron
