# appendix-b-docker/files/rabbitmq.conf
[program:rabbitmq]
command=/docker/rabbitmq.sh
priority=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
