# appendix-b-docker/files/etc/s6-overlay/s6-rc.d/apache2/run
#!/command/execlineb -P
fdmove -c 2 1
apache2ctl -DFOREGROUND
