# appendix-b-docker/files/etc/s6-overlay/s6-rc.d/prepare-koha-envvars/up
/etc/s6-overlay/scripts/01-setup-koha-envvars.sh
