# appendix-b-docker/files/etc/s6-overlay/s6-rc.d/plack/run
#!/command/execlineb -P
fdmove -c 2 1

s6-envdir /etc/koha-envvars

importas instancename INSTANCE_NAME

/usr/bin/starman -M FindBin --socket /var/run/koha/${instancename}/plack.sock --max-requests 50 --workers 2 -E deployment --user=${instancename}-koha --group=${instancename}-koha --access-log /dev/stdout --error-log /dev/stdout "/etc/koha/plack.psgi" 
