# appendix-b-docker/templates/supervisor/cron.conf
[program:cron]
command=/usr/sbin/cron -f -L 5 ; log started and failed cron jobs
priority=100
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0