# appendix-b-docker/templates/supervisor/memcached.conf
[program:memcached]
command=/usr/bin/memcached
priority=10
user=memcache
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
