# appendix-b-docker/templates/supervisor/z3950.conf
[program:z3950]
command= perl /usr/share/koha/bin/z3950_responder.pl
environment=KOHA_CONF="/etc/koha/sites/${KOHA_INSTANCE}/koha-conf.xml",PERL5LIB="/usr/share/koha/lib",KOHA_HOME="/usr/share/koha"
priority=200
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
